# packaging/sshoosh.service
[Unit]
Description=sshoosh SSH workspace chat
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=sshoosh
Group=sshoosh
Environment=SSHOOSH_DB=/var/lib/sshoosh/sshoosh.sqlite
# Environment=SSHOOSH_DATABASE_URL=libsql://example.turso.io
# Environment=SSHOOSH_DATABASE_AUTH_TOKEN=
# Environment=SSHOOSH_NODE_ID=sshoosh-1
# Environment=SSHOOSH_ENCRYPTION_KEY=
Environment=SSHOOSH_MASTER_LEASE_TTL_SECS=15
Environment=SSHOOSH_MASTER_HEARTBEAT_SECS=5
Environment=SSHOOSH_HOST=0.0.0.0
Environment=SSHOOSH_PORT=2222
Environment=SSHOOSH_MAX_CONNECTIONS=256
Environment=SSHOOSH_MAX_CONNECTIONS_PER_IP=32
Environment=SSHOOSH_SERVER_KEY=/var/lib/sshoosh/sshoosh_server_ed25519
ExecStart=/usr/local/bin/sshoosh serve
Restart=on-failure
RestartSec=2
StateDirectory=sshoosh
StateDirectoryMode=0700
UMask=0077
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/var/lib/sshoosh

[Install]
WantedBy=multi-user.target
