# packaging/sshoosh.service
[Unit]
Description=sshoosh SSH workspace chat
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=sshoosh
Group=sshoosh
EnvironmentFile=/etc/sshoosh/sshoosh.env
ExecStart=/usr/local/bin/sshoosh serve
Restart=on-failure
RestartSec=2
StateDirectory=sshoosh
StateDirectoryMode=0700
UMask=0077
NoNewPrivileges=true
CapabilityBoundingSet=
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectClock=true
ProtectHostname=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
RestrictSUIDSGID=true
RestrictNamespaces=true
LockPersonality=true
SystemCallArchitectures=native
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
MemoryMax=256M
TasksMax=128
LimitNOFILE=4096
ReadWritePaths=/var/lib/sshoosh

[Install]
WantedBy=multi-user.target
